# ansible/roles/base/files/start_iptables.timer
[Unit]
Description=Run 1 min after boot
Unit=start_iptables.service

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
